# ansible/roles/hyperglass/files/hyperglass.service
[Unit]
Description=hyperglass
After=network-online.target

[Service]
User=hyperglass
Group=hyperglass
ExecStart=/opt/hyperglass/venv/bin/hyperglass start

[Install]
WantedBy=multi-user.target